# unl0kr/unl0kr-agent.service.in
# SPDX-License-Identifier: GPL-3.0-or-later

[Unit]
Description=Dispatch Password Requests to unl0kr

ConditionPathExists=!/run/plymouth/pid

DefaultDependencies=no
After=plymouth-start.service
Conflicts=emergency.service
Before=emergency.service
Conflicts=shutdown.target initrd-switch-root.target
Before=shutdown.target initrd-switch-root.target

[Service]
ExecStart=@LIBEXECDIR@/unl0kr-agent
